# scripts/server/systemd/fair-crm-backend.service
[Unit]
Description=Fair CRM Backend API
After=network-online.target kyrox-core.service docker.service
Wants=network-online.target
Requires=kyrox-core.service

[Service]
Type=simple
User=@DEPLOY_SERVICE_USER@
Group=@DEPLOY_SERVICE_USER@
WorkingDirectory=@FAIR_CRM_DIR@/backend
Environment=PYTHONPATH=@FAIR_CRM_DIR@/backend
EnvironmentFile=-@FAIR_CRM_DIR@/backend/.env
ExecStart=@FAIR_CRM_DIR@/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
