# scripts/server/systemd/fair-crm-backend.service
[Unit]
Description=Fair CRM Backend API
After=network-online.target kyrox-core.service docker.service
Wants=network-online.target
Requires=kyrox-core.service

[Service]
Type=simple
User=__DEPLOY_SERVICE_USER__
Group=__DEPLOY_SERVICE_USER__
WorkingDirectory=__FAIR_CRM_DIR__/backend
Environment=PYTHONPATH=__FAIR_CRM_DIR__/backend
EnvironmentFile=-__FAIR_CRM_DIR__/backend/.env
ExecStart=__FAIR_CRM_DIR__/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
