# scripts/server/systemd/kyrox-core.service
[Unit]
Description=KYROX Core API
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=@DEPLOY_SERVICE_USER@
Group=@DEPLOY_SERVICE_USER@
WorkingDirectory=@KYROX_CORE_DIR@/backend
Environment=PYTHONPATH=@KYROX_CORE_DIR@/backend
EnvironmentFile=-@KYROX_CORE_DIR@/backend/.env
ExecStart=@KYROX_CORE_DIR@/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
