# scripts/server/systemd/kyrox-core.service
[Unit]
Description=KYROX Core API
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=@DEPLOY_SERVICE_USER@
Group=@DEPLOY_SERVICE_USER@
WorkingDirectory=/opt/kyrox-core/backend
Environment=PYTHONPATH=/opt/kyrox-core/backend
EnvironmentFile=-/opt/kyrox-core/backend/.env
ExecStart=/opt/kyrox-core/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
