# scripts/server/systemd/kyrox-core.service
[Unit]
Description=KYROX Core API
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_SERVICE_USER__
Group=__DEPLOY_SERVICE_USER__
WorkingDirectory=__KYROX_CORE_DIR__/backend
Environment=PYTHONPATH=__KYROX_CORE_DIR__/backend
EnvironmentFile=-__KYROX_CORE_DIR__/backend/.env
ExecStart=__KYROX_CORE_DIR__/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
